# debian/icinga-idoutils.prerm
#!/bin/sh

set -e
#set -x

. /usr/share/debconf/confmodule
. /usr/share/dbconfig-common/dpkg/prerm
dbc_go icinga-idoutils $@

if [ -x "/etc/init.d/ido2db" ]; then
	invoke-rc.d ido2db stop 3>&- || exit $?
fi

#DEBHELPER#
